# infra/refresh/crontab
# Cheap match sync + aggregates + parse requests: every 6 hours. Enqueued —
# the API poller executes; queue order preserves the pipeline sequence.
20 */6 * * * cd /app && ./node_modules/.bin/tsx scripts/enqueue-job.ts fetch-data populate-builds request-parses > /proc/1/fd/1 2>&1
# Slow build fetchers + profile re-sync: daily, off-peak. ~300-500 OpenDota calls.
40 3 * * * cd /app && ./node_modules/.bin/tsx scripts/enqueue-job.ts fetch-hero-builds fetch-player-builds refresh-profiles > /proc/1/fd/1 2>&1
# Nightly DB backup to $BACKUP_DIR (custom-format dump, 7-day rotation).
# Stays a DIRECT run: pg_dump + the backups volume exist only in this image.
10 4 * * * cd /app && ./node_modules/.bin/tsx scripts/run-job.ts backup-db > /proc/1/fd/1 2>&1
